# tools/omnissm/init/systemd/omnissm-inventory.service
[Unit]
Description=Service for OMNISSM Inventory

[Service]
Type=simple
ExecStart=omnissm ps
StandardError=journal
